# supervisord-ipython.conf
[supervisord]
nodaemon=true

[program:ipython]
command=/usr/local/bin/ipython notebook --profile=nbserver --pylab inline
autorestart=true

